# rootfs/etc/s6-overlay/s6-rc.d/filed/run
#!/command/execlineb -P
# shellcheck disable=all

export FILED_STATE_FILE "/var/filed/filed.db"
filed /var/filed
