# rootfs/etc/s6-overlay/s6-rc.d/dufs/run
#!/command/execlineb -P
# shellcheck disable=all

dufs -A /var/filed
